# router/debian/systemd/v4gw-dnsmasq.service
[Unit]
Description=v4-with-v6-nh testbed router dnsmasq (RA + DHCPv4 + return routes)
After=v4gw-router.service
Requires=v4gw-router.service
[Service]
Type=simple
ExecStart=/usr/sbin/dnsmasq --conf-file=/etc/v4gw/dnsmasq.conf --no-daemon
Restart=on-failure
RestartSec=2
[Install]
WantedBy=multi-user.target
