# router/debian/systemd/v4gw-router.service
[Unit]
Description=v4-with-v6-nh testbed router setup (192.0.0.11 gateway, RFC 5549)
After=network.target
Wants=network.target
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/v4gw/setup.sh enp0s2
[Install]
WantedBy=multi-user.target
